# go.ini
[supervisord]
nodaemon=no

[supervisorctl]

[program:gin-live-reload]
command: gin --path src --port 8080 --appPort 8087 run /go/src/app/main.go
