# supervisord.conf
[supervisord]
nodaemon=no

[program:gin-live-reload]
command: gin --path src --port 8080 --appPort 8087 run main.go
